# overlay/var/spool/cron/crontabs/compteur2.sh
#!/bin/sh

echo "Mon test2" >> /tmp/compteur.txt
